# ansible/roles/zabbix/files/zabbix_backup.service
[Unit]
Description=Zabbix Backup Script

[Service]
User=zabbix
ExecStart=/opt/zabbix_backup_scripts/zabbix_backup.sh