# ansible/roles/zabbix/files/zabbix_backup.timer
[Unit]
Description=Run Zabbix Backup Script every 3 days

[Timer]
OnCalendar=*-*-* 12:00:00/3
Persistent=true

[Install]
WantedBy=timers.target